# record-psensors.service
[Unit]
Description=pSensor screenshot service

[Service]
ExecStart=/usr/local/bin/record-sensors.sh

[Install]
WantedBy=graphical.target
